# deploy/systemd/api-center.service
[Unit]
Description=API Center local service
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/scripts/start-local.sh
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=7940
User=__RUN_USER__

[Install]
WantedBy=multi-user.target
